# docker/rootfs/etc/services.d/transmission-telegram-bot/run
#!/usr/bin/with-contenv sh
ttbot -c /opt/transmission-telegram-bot/config.yml
